# zr-check-in/assets/systemd/probe-nexus.service
[Unit]
Description=Probe Nexus Daemon
After=network.target

[Service]
Type=simple
ExecStart=PROBE_BINARY_PATH nexus --wallet WALLET_NAME
Restart=always
RestartSec=10
StandardOutput=append:USER_HOME/.probe/nexus/daemon.log
StandardError=append:USER_HOME/.probe/nexus/daemon.log

[Install]
WantedBy=default.target
