# zr-check-in/assets/systemd/probe-nexus.service
[Unit]
Description=Probe Nexus Daemon - ZENON Red
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=NODE_BIN PROBE_JS nexus --wallet GITHUB_USER
Restart=on-failure
RestartSec=10
Environment=HOME=HOME_DIR
Environment=PATH=NODE_DIR:/usr/local/bin:/usr/bin:/bin
WorkingDirectory=HOME_DIR

[Install]
WantedBy=default.target
